# archive/v1/ubu_workers/airgenome-launcher-cap.timer
[Unit]
Description=AG4 launcher_cap timer — run every 20s

[Timer]
OnBootSec=10s
OnUnitActiveSec=20s
AccuracySec=1s

[Install]
WantedBy=timers.target
