# archive/v1/scripts/com.airgenome.infra-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.infra-probe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/target/release/hexa</string>
        <string>/Users/ghost/Dev/airgenome/modules/infra_probe.hexa</string>
        <string>probe</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/airgenome-infra-probe.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/airgenome-infra-probe.err</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
