# shared/launchagents/com.airgenome.daemons.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 부팅 시 1회 실행 → nohup daemon 으로 compute-tick / remote-load 기동.
         daemon 의 실제 수명은 nohup 자식이 담당. launchd 는 trigger 로만 관여. -->
    <key>Label</key>
    <string>com.airgenome.daemons</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ghost/Dev/airgenome/bin/daemons_start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/.airgenome/daemons.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ghost/.airgenome/daemons.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/airgenome</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
    </dict>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
