# shared/launchagents/com.airgenome.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Airgenome.app/Contents/MacOS/Airgenome</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/.airgenome/menubar.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/.airgenome/menubar.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/airgenome</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
