# shared/launchagents/com.airgenome.stability.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/ghost</string>
		<key>LANG</key>
		<string>en_US.UTF-8</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>Label</key>
	<string>com.airgenome.stability</string>
	<key>LowPriorityIO</key>
	<true/>
	<key>Nice</key>
	<integer>10</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>JetsamPriority</key>
	<integer>40</integer>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/ghost/Dev/airgenome/bin/stability_monitor</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/ghost/.airgenome/stability.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/ghost/.airgenome/stability.stdout.log</string>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>ThrottleInterval</key>
	<integer>60</integer>
	<key>WorkingDirectory</key>
	<string>/Users/ghost/Dev/airgenome</string>
</dict>
</plist>
