# shared/launchagents/com.airgenome.probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.probe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/hexa</string>
        <string>run</string>
        <string>/Users/ghost/Dev/airgenome/modules/probe.hexa</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/.airgenome/probe.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/.airgenome/probe.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/airgenome</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
