# scripts/com.airgenome.ag3-feed-once.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ag3-feed-once — Ubuntu 상태 수집 (ONCE 모드)
  StartInterval=5: launchd가 5초 간격으로 프로세스 1회 실행
  while-true + sleep 누적 CPU 폭주 방지 (D1 골화 2026-04-10)
  기존 com.airgenome.ag3-feed.plist (KeepAlive=true) → 폴백용 보존
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.ag3-feed-once</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.hx/bin/hexa</string>
        <string>__HOME__/Dev/airgenome/modules/ag3_menubar_feed.hexa</string>
        <string>5</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StartInterval</key>
    <integer>5</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>AIRGENOME_AG3FEED_ONCE</key>
        <string>1</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.cargo/bin:__HOME__/.hx/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/airgenome-ag3-feed-once.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/airgenome-ag3-feed-once.err</string>
    <key>ThrottleInterval</key>
    <integer>3</integer>
</dict>
</plist>
