# scripts/com.airgenome.ag3-feed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ag3_menubar_feed — Ubuntu 상태 수집 데몬 (5s 간격 내장)
  내부 while-loop이 5초 간격으로 동작하므로 KeepAlive로 프로세스 유지
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.ag3-feed</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Dev/hexa-lang/target/release/hexa</string>
        <string>__HOME__/Dev/airgenome/modules/ag3_menubar_feed.hexa</string>
        <string>5</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.cargo/bin:__HOME__/.hx/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/airgenome-ag3-feed.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/airgenome-ag3-feed.err</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
