# scripts/com.airgenome.sampler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  sampler — 리소스 모니터 + vitals_ring tick (adaptive 간격 내장)
  내부 while-loop이 3~15초 간격으로 동작하므로 KeepAlive로 프로세스 유지
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.sampler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>export AIRGENOME_SAMPLER_ONCE=1 HOME=/Users/ghost; exec /Users/ghost/.hx/bin/hexa /Users/ghost/Dev/airgenome/sampler.hexa /tmp/airgenome-state.json /Users/ghost/.airgenome/config.json 24576 8 /Users/ghost/Dev/airgenome</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StartInterval</key>
    <integer>15</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>AIRGENOME_SAMPLER_ONCE</key>
        <string>1</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/ghost/.cargo/bin:/Users/ghost/.hx/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/airgenome-sampler.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/airgenome-sampler.err</string>
    <key>ThrottleInterval</key>
    <integer>3</integer>
</dict>
</plist>
