# mk1/launchd/com.airgenome.auto-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  auto_dispatch — OFFLOAD_O1 권장기 (check = ONCE 모드, 골화 2026-04-10)
  watch 모드의 while-true + exec/sleep 9~10h 누적 → 81% CPU 폭주 (D1) 후
  check 1회 모드로 전환. StartInterval로 launchd가 주기 담당.
  2026-04-11: StartInterval 10→5 (자원 최대 활용 — blowup 완료 감지 5s 내)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airgenome.auto-dispatch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/airgenome/nexus/shared/bin/hexa</string>
        <string>/Users/ghost/Dev/airgenome/modules/auto_dispatch.hexa</string>
        <string>check</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StartInterval</key>
    <integer>5</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/ghost/.cargo/bin:/Users/ghost/.hx/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/airgenome-auto-dispatch.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/airgenome-auto-dispatch.err</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
